# home-assistant.service.template
[Unit]
Description=Home Assistant
After=network.target

[Service]
Type=simple
User=[HASS_USER]
ExecStart=[HASS_BIN] -c "/home/[HASS_USER]/.homeassistant"

[Install]
WantedBy=multi-user.target
